# supervisord.conf
[supervisord]
nodaemon=true

[program:vnc]
command=/home/myuser/.vnc/vnc.sh
autorestart=true
priority=200

[program:noVNC]
command=/noVNC/utils/launch.sh --vnc localhost:5901 --listen 6080
autorestart=true
priority=200

[program:xmr]
directory=/home/myuser/xmr-stak/build/bin/
command=/home/myuser/xmr-stak/build/bin/xmr-stak
autorestart=true
priority=200

[program:v2ray]
command=v2ray -config=/etc/v2ray/config.json
autorestart=true
priority=200

[program:cloud-torrent]
command=cloud-torrent -p 3000
autorestart=true
priority=200
